# k8s/aws/packer/files/fpga.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/fgpa-init.sh

[Install]
WantedBy=default.target